# util/udev/60-onbattery.rules
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="0",RUN+="/usr/bin/profile-battery.sh"
