# util/udev/61-onpower.rules
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/bin/profile-power.sh"
